# deploy/config/common.local/etc/uv4l.service
[Unit]
Description=uv4l daemon service

[Service]
ExecStart=/usr/local/bin/run-uv4l.sh
SyslogIdentifier=uv4l

[Install]
WantedBy=multi-user.target

